# scripts/themealdb/com.retsepti.themealdb-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  macOS launchd job — runs the TheMealDB sync once per day at 04:00 local time.

  Install:
    1) Edit "WORKING_DIR" below if your repo path differs.
    2) cp scripts/themealdb/com.retsepti.themealdb-sync.plist ~/Library/LaunchAgents/
    3) launchctl load ~/Library/LaunchAgents/com.retsepti.themealdb-sync.plist
    4) (optional) launchctl start com.retsepti.themealdb-sync

  Logs: <repo>/.cache/themealdb/sync.log
  Uninstall:
    launchctl unload ~/Library/LaunchAgents/com.retsepti.themealdb-sync.plist
    rm ~/Library/LaunchAgents/com.retsepti.themealdb-sync.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.retsepti.themealdb-sync</string>

    <key>WorkingDirectory</key>
    <string>/Users/abrosnahat/Desktop/retsepti</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>npm run sync:themealdb</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/abrosnahat/Desktop/retsepti/.cache/themealdb/sync.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/abrosnahat/Desktop/retsepti/.cache/themealdb/sync.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
